# files/etc/uci-defaults/96-install-extras
#!/bin/sh
# 96-install-extras
# 首次启动后异步安装所有移出固件的软件包
# 包括：sing-box / homeproxy LuCI / dockerd / dockerman LuCI /
#       adguardhome / samba4 / lucky / smartmontools / bash / htop 等

DONE_FLAG="/etc/install-extras-done"
[ -f "$DONE_FLAG" ] && exit 0

# 写入 procd init 服务（传统 rc.common，start-stop-daemon 后台 fork）
cat > /etc/init.d/install-extras << 'INITEOF'
#!/bin/sh /etc/rc.common
START=99
STOP=10

start() {
    [ -f /etc/install-extras-done ] && return 0
    [ -f /usr/lib/install-extras/run.sh ] || return 1
    start-stop-daemon -S -b -m -p /var/run/install-extras.pid \
        -x /bin/sh -- /usr/lib/install-extras/run.sh
}

stop() {
    start-stop-daemon -K -p /var/run/install-extras.pid 2>/dev/null || true
}
INITEOF
chmod +x /etc/init.d/install-extras

mkdir -p /usr/lib/install-extras
cat > /usr/lib/install-extras/run.sh << 'RUNEOF'
#!/bin/sh
LOG_TAG="install-extras"
DONE_FLAG="/etc/install-extras-done"

log() { echo "[$(date '+%H:%M:%S')] $1" | logger -t "$LOG_TAG"; }

[ -f "$DONE_FLAG" ] && exit 0

log "⏳ 等待系统启动完成 (60s)..."
sleep 60

# 修复 apk video feed
for f in /etc/apk/repositories.d/*.list; do
    [ -f "$f" ] && sed -i '/\/video\//d' "$f"
done

# 等待网络（最多 8 分钟）
log "⏳ 等待网络就绪..."
NETWORK_OK=0
for i in $(seq 1 240); do
    ping -c 1 -W 2 223.5.5.5 >/dev/null 2>&1 \
        || ping -c 1 -W 2 8.8.8.8 >/dev/null 2>&1 \
        && { NETWORK_OK=1; log "✅ 网络就绪（第 ${i} 次）"; break; }
    sleep 2
done

if [ "$NETWORK_OK" -eq 0 ]; then
    log "❌ 网络超时，跳过本次安装，下次启动自动重试"
    /etc/init.d/install-extras disable 2>/dev/null || true
    exit 1
fi

apk update 2>&1 | logger -t "$LOG_TAG"

# ── 安装列表 ──────────────────────────────────────────────────────────
# 格式: "包名1 包名2 ..."  每组作为一条 apk add 命令
PKGS="
sing-box
luci-app-homeproxy luci-i18n-homeproxy-zh-cn
luci-app-dockerman luci-i18n-dockerman-zh-cn
dockerd docker docker-compose cgroupfs-mount
adguardhome luci-app-adguardhome luci-i18n-adguardhome-zh-cn
samba4-server luci-app-samba4 luci-i18n-samba4-zh-cn
lucky luci-app-lucky luci-i18n-lucky-zh-cn
luci-app-fileassistant
smartmontools
bash htop lsblk fdisk
luci-app-crontabs luci-i18n-crontabs-zh-cn
luci-app-wol luci-i18n-wol-zh-cn etherwake
"

echo "$PKGS" | while IFS= read -r line; do
    [ -z "$(echo "$line" | tr -d ' ')" ] && continue
    log "📦 安装: $line"
    # shellcheck disable=SC2086
    apk add --no-cache $line 2>&1 | logger -t "$LOG_TAG" \
        && log "✅ 完成: $line" \
        || log "⚠️  失败: $line（可手动: apk add $line）"
done

# ── sing-box 安装后重启 homeproxy ─────────────────────────────────────
if command -v sing-box >/dev/null 2>&1; then
    /etc/init.d/homeproxy restart 2>/dev/null || true
    log "✅ homeproxy 已重启"
fi

# ── 配置 Docker ───────────────────────────────────────────────────────
if command -v dockerd >/dev/null 2>&1; then
    mkdir -p /etc/docker
    USB_MOUNT=""
    for mp in /mnt/sda1 /mnt/sda /mnt/usb /mnt/mmcblk0p1; do
        mountpoint -q "$mp" 2>/dev/null && { USB_MOUNT="$mp"; break; }
    done
    if [ -n "$USB_MOUNT" ]; then
        DOCKER_DIR="$USB_MOUNT/docker-data"
        mkdir -p "$DOCKER_DIR"
        printf '{\n  "data-root": "%s",\n  "log-driver": "json-file",\n  "log-opts": {"max-size": "10m", "max-file": "3"}\n}\n' \
            "$DOCKER_DIR" > /etc/docker/daemon.json
        log "✅ Docker 数据目录: $DOCKER_DIR (USB)"
    else
        printf '{\n  "log-driver": "json-file",\n  "log-opts": {"max-size": "5m", "max-file": "2"}\n}\n' \
            > /etc/docker/daemon.json
        log "⚠️  Docker 使用 overlay（建议插 USB）"
    fi
    mountpoint -q /sys/fs/cgroup || \
        mount -t cgroup2 none /sys/fs/cgroup 2>/dev/null || true
    /etc/init.d/dockerd enable 2>/dev/null || true
    /etc/init.d/dockerd start  2>/dev/null || true
    sleep 3
    /etc/init.d/dockerd status 2>/dev/null | grep -q running \
        && log "✅ dockerd 已启动" \
        || log "⚠️  dockerd 未运行，手动: /etc/init.d/dockerd start"
fi

# ── 刷新 LuCI ─────────────────────────────────────────────────────────
rm -f /tmp/luci-indexcache* 2>/dev/null
/etc/init.d/rpcd restart 2>/dev/null || true

/etc/init.d/install-extras disable 2>/dev/null || true
touch "$DONE_FLAG"
log "🎉 install-extras 全部完成"
RUNEOF

chmod +x /usr/lib/install-extras/run.sh
/etc/init.d/install-extras enable 2>/dev/null || true
logger -t "96-install-extras" "✅ install-extras 服务已注册"
exit 0
