# files/etc/uci-defaults/96-install-extras
#!/bin/sh
# 96-install-extras [v4.3 - Final]
# 修复：锁文件 PID 写法 / trap 时序 / DNS 恢复空值 / 自愈重试上限 / block hotplug 防重复
# 适用：360 V6 / IPQ6000 (ImmortalWrt / ash)

DONE_FLAG="/etc/install-extras-done"
USB_DONE="/etc/install-usb-done"
STATE_DIR="/etc/install-state"
MOUNT_FILE="/etc/install-usb-mount"

mkdir -p "$STATE_DIR" /usr/lib/install-extras /usr/bin \
         /etc/hotplug.d/iface /etc/hotplug.d/block \
         /etc/dbus-1 /etc/profile.d /etc/crontabs

_is_mounted() { grep -q "^[^ ]* $1 " /proc/mounts 2>/dev/null; }

# ====================================================================
# 模块 0：状态查询
# ====================================================================
cat > /usr/bin/360v6-status << 'STATUSEOF'
#!/bin/sh
echo "🔍 状态 @ $(date '+%H:%M:%S')"
[ -f /etc/install-extras-done ] && echo "✅ 基础：完成" || echo "⏳ 基础：待触发"
[ -f /etc/install-usb-done ]    && echo "✅ 重型：完成" || echo "⏳ 重型：待 U 盘"
echo "💾 Flash: $(df -h /overlay 2>/dev/null | awk 'NR==2{print $5}')"
docker info >/dev/null 2>&1    && echo "🐳 Docker: OK" || echo "🐳 Docker: --"
pgrep AdGuardHome >/dev/null 2>&1 && echo "🛡️  AGH: OK"    || echo "🛡️  AGH: --"
cat /etc/install-usb-mount 2>/dev/null | xargs -I{} echo "📂 U盘: {}"
p1=$(cat /etc/install-state/phase1-failures 2>/dev/null || echo 0)
p2=$(cat /etc/install-state/phase2-failures 2>/dev/null || echo 0)
[ "$p1" -gt 0 ] && echo "⚠️  基础失败次数: $p1/3"
[ "$p2" -gt 0 ] && echo "⚠️  重型失败次数: $p2/3"
echo ""
echo "手动重跑基础: rm -f /etc/install-extras-done /etc/install-state/phase1-failures && /usr/lib/install-extras/run.sh"
echo "手动重跑重型: rm -f /etc/install-usb-done   /etc/install-state/phase2-failures && /usr/lib/install-extras/install-usb.sh"
STATUSEOF
chmod +x /usr/bin/360v6-status

# ====================================================================
# 模块 A：install-usb.sh（核心部署 + 自愈）
# ====================================================================
cat > /usr/lib/install-extras/install-usb.sh << 'USBEOF'
#!/bin/sh
LOG_TAG="usb"
USB_DONE="/etc/install-usb-done"
LOCK="/tmp/install-usb.lock"
STATE_DIR="/etc/install-state"
MOUNT_FILE="/etc/install-usb-mount"

log() { echo "[$(date '+%H:%M:%S')] $1" | logger -t "$LOG_TAG"; }
_is_mounted() { grep -q "^[^ ]* $1 " /proc/mounts 2>/dev/null; }

mkdir -p "$STATE_DIR"

# ── 自愈模式（--self-heal）────────────────────────────────────────
# Fix #1：先检查重试上限，再检查锁进程存活，避免无限重启
if [ "${1:-}" = "--self-heal" ]; then
    [ -f "$USB_DONE" ] && exit 0
    __sf=$(cat "$STATE_DIR/phase2-failures" 2>/dev/null || echo 0)
    [ "${__sf:-0}" -ge 3 ] && exit 0   # 失败过多，停止自愈
    if [ -f "$LOCK" ]; then
        _pid=$(cat "$LOCK" 2>/dev/null)
        kill -0 "$_pid" 2>/dev/null && exit 0   # 进程活着，正常跳过
        log "🔧 清理残留锁 (旧pid=$_pid)"
        rm -f "$LOCK"
    fi
    for _mp in /mnt/sda1 /mnt/sda /mnt/sdb1 /mnt/sdb /mnt/usb /mnt/mmcblk0p1; do
        _is_mounted "$_mp" && {
            log "🔧 自愈触发: $_mp"
            /usr/lib/install-extras/install-usb.sh "$_mp" &
            exit 0
        }
    done
    exit 0
fi

# ── 幂等 / 锁保护 ────────────────────────────────────────────────
[ -f "$USB_DONE" ] && { log "ℹ️ 已完成"; exit 0; }

# Fix #2：先设 trap，再写锁文件，防止写锁和 trap 之间被 kill 导致锁残留
# Fix #3：echo $$ 而非 echo $$$$（后者输出双倍 PID，如 12341234）
trap 'rm -f "$LOCK"' EXIT INT TERM
if [ -f "$LOCK" ]; then
    _pid=$(cat "$LOCK" 2>/dev/null)
    kill -0 "$_pid" 2>/dev/null && { log "⚠️ 运行中(pid=$_pid)"; exit 0; }
    log "🔧 清理残留锁 (旧pid=$_pid)"
    rm -f "$LOCK"
fi
echo $$ > "$LOCK"

# 重试上限
__fail=$(cat "$STATE_DIR/phase2-failures" 2>/dev/null || echo 0)
[ "${__fail:-0}" -ge 3 ] && { log "❌ 失败次数过多，停止重试"; exit 1; }

# ── 确定 U 盘挂载点 ──────────────────────────────────────────────
USB_MOUNT="${1:-}"
case "$USB_MOUNT" in --*) USB_MOUNT="" ;; esac
[ -z "$USB_MOUNT" ] && {
    for _mp in /mnt/sda1 /mnt/sda /mnt/sdb1 /mnt/sdb /mnt/usb /mnt/mmcblk0p1; do
        _is_mounted "$_mp" && { USB_MOUNT="$_mp"; break; }
    done
}
[ -z "$USB_MOUNT" ] && \
    USB_MOUNT=$(grep -E '^/dev/sd[a-z][0-9]* ' /proc/mounts 2>/dev/null | head -1 | awk '{print $2}')
[ -z "$USB_MOUNT" ] && { log "❌ 无 U 盘挂载点"; echo $((__fail + 1)) > "$STATE_DIR/phase2-failures"; exit 1; }
echo "$USB_MOUNT" > "$MOUNT_FILE"
log "▶ 部署 | U盘: $USB_MOUNT"

# ── TMPDIR / apk 缓存重定向到 U 盘 ──────────────────────────────
export TMPDIR="$USB_MOUNT/.tmp"
APK_CACHE="$USB_MOUNT/.apk-cache"
mkdir -p "$TMPDIR" "$APK_CACHE"

# ── DNS 保护 ─────────────────────────────────────────────────────
# Fix #4：__dns_old 若为空则写 nameserver 127.0.0.1 而非空文件
__dns_old=$(cat /etc/resolv.conf 2>/dev/null)
[ -z "$__dns_old" ] && __dns_old="nameserver 127.0.0.1"
printf 'nameserver 223.5.5.5\nnameserver 8.8.8.8\n' > /etc/resolv.conf
# 覆盖 trap，加入 DNS 恢复（此时 __dns_old 已有安全值）
trap 'rm -f "$LOCK"; printf "%s\n" "$__dns_old" > /etc/resolv.conf; /etc/init.d/dnsmasq restart >/dev/null 2>&1' EXIT INT TERM

# ── 空间预检（≥ 2GB）─────────────────────────────────────────────
_avail=$(df -k "$USB_MOUNT" 2>/dev/null | awk 'NR==2{print $4}')
if [ "${_avail:-0}" -lt 2097152 ]; then
    log "⚠️ U 盘空间不足 (${_avail:-0}KB < 2GB)，清理缓存后退出"
    rm -rf "$USB_MOUNT/.apk-cache" "$USB_MOUNT/.tmp" 2>/dev/null
    echo $((__fail + 1)) > "$STATE_DIR/phase2-failures"
    exit 1
fi
echo "10" > /proc/sys/vm/swappiness 2>/dev/null

# ── Swap ─────────────────────────────────────────────────────────
SWAPFILE="$USB_MOUNT/swapfile"
rm -f "${SWAPFILE}.tmp" 2>/dev/null
if ! grep -q swapfile /proc/swaps 2>/dev/null; then
    if [ ! -f "$SWAPFILE" ]; then
        dd if=/dev/zero of="${SWAPFILE}.tmp" bs=1M count=256 2>/dev/null && {
            mv "${SWAPFILE}.tmp" "$SWAPFILE"
            sync; chmod 600 "$SWAPFILE"
            mkswap "$SWAPFILE" >/dev/null 2>&1
            log "✅ Swap 已创建"
        } || rm -f "${SWAPFILE}.tmp"
    fi
    [ -f "$SWAPFILE" ] && swapon "$SWAPFILE" 2>/dev/null && {
        grep -qF "$SWAPFILE" /etc/fstab 2>/dev/null || echo "$SWAPFILE none swap sw 0 0" >> /etc/fstab
        log "✅ Swap 已激活"
    }
fi

# ── 等待 apk 就绪 ────────────────────────────────────────────────
_i=1
while [ "$_i" -le 30 ]; do
    apk info >/dev/null 2>&1 && break
    sleep 2; _i=$((_i + 1))
done
apk info >/dev/null 2>&1 || { log "❌ apk 长时间锁定"; echo $((__fail + 1)) > "$STATE_DIR/phase2-failures"; exit 1; }

# ── 安装函数 ─────────────────────────────────────────────────────
__ifail=0
__install() {
    _t=$(mktemp "$TMPDIR/apk.XXXXXX")
    log "📦 $*"
    apk add --cache-dir "$APK_CACHE" "$@" > "$_t" 2>&1
    _r=$?; cat "$_t" | logger -t "$LOG_TAG"; rm -f "$_t"
    [ "$_r" -ne 0 ] && { log "❌ 失败: $*"; __ifail=$((__ifail + 1)); }
    return $_r
}

# ── Docker ───────────────────────────────────────────────────────
DS="$USB_MOUNT/docker-system"
DD="$USB_MOUNT/docker-data"
mkdir -p "$DS" "$DD" /etc/docker
if [ ! -f "$DS/etc/apk/world" ]; then
    mkdir -p "$DS/etc/apk" "$DS/var/cache/apk"
    apk add --root "$DS" --initdb --no-cache 2>/dev/null || true
fi
__install --root "$DS" docker dockerd

__dok=0
if [ -f "$DS/usr/bin/dockerd" ]; then
    mkdir -p /overlay/upper/usr/bin
    for _b in docker dockerd containerd runc; do
        ln -sf "$DS/usr/bin/$_b" "/overlay/upper/usr/bin/$_b" 2>/dev/null || true
    done
    echo "export LD_LIBRARY_PATH=\"$DS/usr/lib:\$LD_LIBRARY_PATH\"" > /etc/profile.d/docker-usb.sh
    cat > /etc/docker/daemon.json << DJSON
{"data-root":"$DD","storage-driver":"overlay2","log-driver":"json-file","log-opts":{"max-size":"5m","max-file":"2"}}
DJSON
    cat > /etc/init.d/dockerd << 'DINIT'
#!/bin/sh /etc/rc.common
START=99; STOP=10; USE_PROCD=1
start_service() {
    [ -f /etc/install-usb-done ] || return 0
    M=$(cat /etc/install-usb-mount 2>/dev/null)
    [ -n "$M" ] || return 1
    [ -f "$M/docker-system/usr/bin/dockerd" ] || return 1
    export LD_LIBRARY_PATH="$M/docker-system/usr/lib:$LD_LIBRARY_PATH"
    procd_open_instance
    procd_set_param command "$M/docker-system/usr/bin/dockerd" --config-file /etc/docker/daemon.json
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param respawn
    procd_close_instance
}
stop_service() { pkill dockerd 2>/dev/null; pkill containerd 2>/dev/null; return 0; }
DINIT
    chmod +x /etc/init.d/dockerd
    /etc/init.d/dockerd enable  2>&1 | logger -t "$LOG_TAG"
    /etc/init.d/dockerd start   2>&1 | logger -t "$LOG_TAG"
    _i=1
    while [ "$_i" -le 15 ]; do
        docker info >/dev/null 2>&1 && { __dok=1; break; }
        sleep 2; _i=$((_i + 1))
    done
    [ "$__dok" -eq 1 ] && log "✅ Docker 已启动" || log "⚠️ Docker 启动超时（稍后可能就绪）"
fi

# ── AdGuardHome ──────────────────────────────────────────────────
AS="$USB_MOUNT/adguardhome-system"
AD="$USB_MOUNT/adguardhome-data"
AC="$USB_MOUNT/adguardhome-config"
mkdir -p "$AS" "$AD" "$AC"
if [ ! -f "$AS/etc/apk/world" ]; then
    mkdir -p "$AS/etc/apk"
    apk add --root "$AS" --initdb --no-cache 2>/dev/null || true
fi
__install --root "$AS" adguardhome

__aok=0
if [ -f "$AS/usr/bin/AdGuardHome" ]; then
    mkdir -p /overlay/upper/usr/bin
    ln -sf "$AS/usr/bin/AdGuardHome" "/overlay/upper/usr/bin/AdGuardHome" 2>/dev/null || true

    [ ! -f "$AC/AdGuardHome.yaml" ] && [ -f "/etc/adguardhome/AdGuardHome.yaml" ] && \
        cp -a /etc/adguardhome/. "$AC/" 2>/dev/null
    [ -z "$(ls -A "$AD" 2>/dev/null)" ] && [ -d "/var/lib/adguardhome" ] && \
        cp -a /var/lib/adguardhome/. "$AD/" 2>/dev/null

    uci set adguardhome.config.config_file="$AC/AdGuardHome.yaml" 2>/dev/null || true
    uci set adguardhome.config.work_dir="$AD" 2>/dev/null || true
    uci commit adguardhome 2>/dev/null || true

    cat > /etc/init.d/adguardhome << 'AINIT'
#!/bin/sh /etc/rc.common
START=98; STOP=10; USE_PROCD=1
start_service() {
    [ -f /etc/install-usb-done ] || return 0
    M=$(cat /etc/install-usb-mount 2>/dev/null)
    [ -n "$M" ] || return 1
    [ -f "$M/adguardhome-system/usr/bin/AdGuardHome" ] || return 1
    procd_open_instance
    procd_set_param command "$M/adguardhome-system/usr/bin/AdGuardHome" \
        --config "$M/adguardhome-config/AdGuardHome.yaml" \
        --work-dir "$M/adguardhome-data"
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param respawn
    procd_close_instance
}
stop_service() { pkill AdGuardHome 2>/dev/null; return 0; }
AINIT
    chmod +x /etc/init.d/adguardhome
    /etc/init.d/adguardhome enable  >/dev/null 2>&1
    /etc/init.d/adguardhome restart >/dev/null 2>&1
    sleep 3
    pgrep -x AdGuardHome >/dev/null 2>&1 && __aok=1
    [ "$__aok" -eq 1 ] && log "✅ AGH 已启动" || log "⚠️ AGH 进程未就绪"
fi

# ── 刷新 LuCI ────────────────────────────────────────────────────
rm -f /tmp/luci-indexcache* 2>/dev/null
/etc/init.d/rpcd restart >/dev/null 2>&1

# ── 收尾 ─────────────────────────────────────────────────────────
if [ "$__dok" -eq 1 ] || [ "$__aok" -eq 1 ]; then
    touch "$USB_DONE"
    rm -f "$STATE_DIR/phase2-failures"
    _ds=$([ "$__dok" -eq 1 ] && echo "OK" || echo "FL")
    _as=$([ "$__aok" -eq 1 ] && echo "OK" || echo "FL")
    log "🎉 完成 [Docker:$_ds AGH:$_as]"
    exit 0
fi
echo $((__fail + 1)) > "$STATE_DIR/phase2-failures"
log "❌ 全部失败，等待自愈重试 ($((__fail + 1))/3)"
exit 1
USBEOF
chmod +x /usr/lib/install-extras/install-usb.sh

# ====================================================================
# 模块 B：run.sh（Phase 1 基础环境）
# ====================================================================
cat > /usr/lib/install-extras/run.sh << 'RUNEOF'
#!/bin/sh
LOG_TAG="main"
DONE="/etc/install-extras-done"
USB_DONE="/etc/install-usb-done"
LOCK="/tmp/install-main.lock"
STATE_DIR="/etc/install-state"

log() { echo "[$(date '+%H:%M:%S')] $1" | logger -t "$LOG_TAG"; }
_is_mounted() { grep -q "^[^ ]* $1 " /proc/mounts 2>/dev/null; }

mkdir -p "$STATE_DIR"
[ -f "$DONE" ] && { log "ℹ️ 基础已完成，跳过"; exit 0; }

# Fix #2：先设 trap，再写锁文件
trap 'rm -f "$LOCK"' EXIT INT TERM
if [ -f "$LOCK" ]; then
    _pid=$(cat "$LOCK" 2>/dev/null)
    kill -0 "$_pid" 2>/dev/null && { log "⚠️ 运行中(pid=$_pid)"; exit 0; }
    log "🔧 清理残留锁 (旧pid=$_pid)"
    rm -f "$LOCK"
fi
# Fix #3：echo $$ 写入正确 PID
echo $$ > "$LOCK"

__fail=$(cat "$STATE_DIR/phase1-failures" 2>/dev/null || echo 0)
[ "${__fail:-0}" -ge 3 ] && { log "❌ 失败过多，暂停自动重试"; exit 1; }

log "▶ v4.3 基础环境启动"
sleep 10

# 等待网络（最长 3 分钟）
NET_OK=0; _i=1
while [ "$_i" -le 36 ]; do
    ping -c1 -W2 223.5.5.5 >/dev/null 2>&1 && { NET_OK=1; break; }
    ping -c1 -W2 8.8.8.8   >/dev/null 2>&1 && { NET_OK=1; break; }
    sleep 5; _i=$((_i + 1))
done
[ "$NET_OK" -eq 0 ] && { log "❌ 网络不可达"; echo $((__fail + 1)) > "$STATE_DIR/phase1-failures"; exit 1; }
log "✅ 网络就绪"

# DNS 注入
# Fix #4：__dns_old 为空时保底写 nameserver 127.0.0.1
__dns_old=$(cat /etc/resolv.conf 2>/dev/null)
[ -z "$__dns_old" ] && __dns_old="nameserver 127.0.0.1"
printf 'nameserver 223.5.5.5\nnameserver 8.8.8.8\n' > /etc/resolv.conf
# 覆盖 trap，加入 DNS 恢复（此时 __dns_old 已有安全值）
trap 'rm -f "$LOCK"; printf "%s\n" "$__dns_old" > /etc/resolv.conf; /etc/init.d/dnsmasq restart >/dev/null 2>&1; log "✅ DNS 已恢复"' EXIT INT TERM
log "✅ 临时 DNS 已注入"

# Flash 预检
FLASH=$(df /overlay 2>/dev/null | awk 'NR==2{print $4}')
[ "${FLASH:-0}" -lt 20480 ] && { log "❌ Flash 不足 20MB"; echo $((__fail + 1)) > "$STATE_DIR/phase1-failures"; exit 1; }
log "📊 Flash 剩余: ${FLASH}KB"

# apk update（HTTPS → HTTP 降级）
__apk_update() {
    _t=$(mktemp /tmp/apk.XXXXXX)
    apk update > "$_t" 2>&1
    cat "$_t" | logger -t "$LOG_TAG"
    _n=$(grep -oE 'OK: [0-9]+ distinct' "$_t" 2>/dev/null | grep -oE '[0-9]+' | head -1)
    rm -f "$_t"
    [ "${_n:-0}" -gt 800 ] && return 0
    log "⚠️ HTTPS 包数不足，切换 HTTP"
    for _f in /etc/apk/repositories.d/*.list; do
        [ -f "$_f" ] && sed 's|https://|http://|g' "$_f" > "${_f}.tmp" && mv "${_f}.tmp" "$_f"
    done 2>/dev/null
    _t=$(mktemp /tmp/apk.XXXXXX)
    apk update > "$_t" 2>&1
    cat "$_t" | logger -t "$LOG_TAG"
    _n=$(grep -oE 'OK: [0-9]+ distinct' "$_t" 2>/dev/null | grep -oE '[0-9]+' | head -1)
    rm -f "$_t"
    [ "${_n:-0}" -gt 800 ]
}
__apk_update || { log "❌ apk update 失败"; echo $((__fail + 1)) > "$STATE_DIR/phase1-failures"; exit 1; }
log "✅ apk update 成功"

# 安装函数
__inst_fail=0
__inst() {
    _t=$(mktemp /tmp/apk.XXXXXX)
    log "📦 $*"
    apk add --no-cache "$@" > "$_t" 2>&1
    _r=$?; cat "$_t" | logger -t "$LOG_TAG"; rm -f "$_t"
    [ "$_r" -ne 0 ] && { log "❌ 失败: $*"; __inst_fail=$((__inst_fail + 1)); }
    return $_r
}

__inst bash htop lsblk smartmontools fdisk
[ "${FLASH:-0}" -ge 25000 ] && __inst samba4-server luci-app-samba4 luci-i18n-samba4-zh-cn \
    || log "ℹ️ Flash 不足 25MB，跳过 Samba"

# 修复 dbus/avahi
/etc/init.d/avahi-daemon stop >/dev/null 2>&1; sleep 1
/etc/init.d/dbus stop  >/dev/null 2>&1; sleep 1
/etc/init.d/dbus start >/dev/null 2>&1; sleep 2
/etc/init.d/avahi-daemon start >/dev/null 2>&1
log "✅ dbus/avahi 已重载"

rm -f /tmp/luci-indexcache* 2>/dev/null
/etc/init.d/rpcd restart >/dev/null 2>&1

if [ "$__inst_fail" -ne 0 ]; then
    log "⚠️ $__inst_fail 组包失败，下次 WAN UP 重试"
    echo $((__fail + 1)) > "$STATE_DIR/phase1-failures"
    exit 1
fi

touch "$DONE"
rm -f "$STATE_DIR/phase1-failures"
log "🎉 基础环境完成"

# Phase 2：一次性检查 U 盘（不轮询，避免与系统 hotplug 竞态）
# 已挂载 → 直接触发；未挂载 → 交由 block hotplug 哨兵处理
[ ! -f "$USB_DONE" ] && {
    MP=""
    for _m in /mnt/sda1 /mnt/sda /mnt/sdb1 /mnt/sdb /mnt/usb /mnt/mmcblk0p1; do
        _is_mounted "$_m" && { MP="$_m"; break; }
    done
    [ -z "$MP" ] && MP=$(grep -E '^/dev/sd[a-z][0-9]* ' /proc/mounts 2>/dev/null | head -1 | awk '{print $2}')
    if [ -n "$MP" ]; then
        log "✅ U 盘已挂载: $MP，触发重型应用..."
        /usr/lib/install-extras/install-usb.sh "$MP" &
    else
        log "ℹ️ U 盘未挂载，等待 hotplug 哨兵触发"
    fi
}
RUNEOF
chmod +x /usr/lib/install-extras/run.sh

# ====================================================================
# 模块 C：热插拔触发器
# ====================================================================
cat > /etc/hotplug.d/iface/96-install-extras << 'IFACEEOF'
#!/bin/sh
[ "$ACTION" = "ifup" ] || exit 0
case "$INTERFACE" in wan|wan6|pppoe-wan) ;; *) exit 0 ;; esac
[ -f /etc/install-extras-done ] && exit 0
logger -t "hotplug" "🌐 WAN 上线，触发基础环境"
/usr/lib/install-extras/run.sh &
IFACEEOF
chmod +x /etc/hotplug.d/iface/96-install-extras

# Fix #5：block hotplug 加入防重复进程检查
cat > /etc/hotplug.d/block/96-install-usb << 'BLOCKEOF'
#!/bin/sh
_is_mounted() { grep -q "^[^ ]* $1 " /proc/mounts 2>/dev/null; }
[ "$ACTION" = "add" ]           || exit 0
[ -f /etc/install-usb-done ]    && exit 0
[ -f /etc/install-extras-done ] || exit 0
logger -t "hotplug" "💾 U 盘插入: DEVNAME=${DEVNAME:-?}"
sleep 5
MP=""
for mp in /mnt/sda1 /mnt/sda /mnt/sdb1 /mnt/sdb /mnt/usb /mnt/mmcblk0p1; do
    _is_mounted "$mp" && { MP="$mp"; break; }
done
[ -z "$MP" ] && MP=$(grep -E '^/dev/sd[a-z][0-9]* ' /proc/mounts 2>/dev/null | head -1 | awk '{print $2}')
[ -z "$MP" ] && { logger -t "hotplug" "⚠️ 挂载点未就绪，跳过"; exit 1; }
# 防止重复触发
pgrep -f "install-usb.sh" >/dev/null 2>&1 && {
    logger -t "hotplug" "⚠️ 安装脚本已在运行，跳过"
    exit 0
}
logger -t "hotplug" "💾 U 盘就绪: $MP，触发重型应用"
/usr/lib/install-extras/install-usb.sh "$MP" &
BLOCKEOF
chmod +x /etc/hotplug.d/block/96-install-usb

# ====================================================================
# 模块 D：自愈定时任务（每 5 分钟检查一次）
# ====================================================================
if ! grep -q "self-heal" /etc/crontabs/root 2>/dev/null; then
    echo "*/5 * * * * PATH=/usr/bin:/bin:/usr/sbin:/sbin /usr/lib/install-extras/install-usb.sh --self-heal 2>&1 | logger -t self-heal" >> /etc/crontabs/root
    /etc/init.d/cron restart 2>/dev/null || true
fi

# ====================================================================
# 首次触发（脚本执行时如基础环境未完成则立即启动）
# ====================================================================
[ ! -f /etc/install-extras-done ] && /usr/lib/install-extras/run.sh &

logger -t "install-extras" "✅ v4.3 Final 就绪"
echo "━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━"
echo "🚀 v4.3 Final 部署完成"
echo "📊 状态: 360v6-status"
echo "📋 日志: logread -f | grep -E 'main|usb|hotplug'"
echo "━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━"
exit 0
