# files/etc/uci-defaults/96-install-extras
#!/bin/sh
# 96-install-extras [Final v3.4]

DONE_FLAG="/etc/install-extras-done"
[ -f "$DONE_FLAG" ] && exit 0

mkdir -p /usr/lib/install-extras

# ══════════════════════════════════════════════════════════════
# install-usb.sh：重型应用安装（Docker + AdGuardHome）
# 由 WAN UP 触发（有 U 盘时）或 USB 热插触发器调用
# ══════════════════════════════════════════════════════════════
cat > /usr/lib/install-extras/install-usb.sh << 'USBEOF'
#!/bin/sh
LOG_TAG="install-extras"
USB_DONE="/etc/install-usb-done"
USB_LOCK="/tmp/install-usb.lock"

log() { echo "[$(date '+%H:%M:%S')] $1" | logger -t "$LOG_TAG"; }

[ -f "$USB_DONE" ] && exit 0
[ -f "$USB_LOCK" ] && exit 0
touch "$USB_LOCK"

USB_MOUNT="$1"
# 未传参时自动探测
if [ -z "$USB_MOUNT" ]; then
    for mp in /mnt/sda1 /mnt/sda /mnt/usb /mnt/mmcblk0p1; do
        mountpoint -q "$mp" 2>/dev/null && { USB_MOUNT="$mp"; break; }
    done
fi

if [ -z "$USB_MOUNT" ]; then
    log "❌ install-usb: 未找到 U 盘挂载点"
    rm -f "$USB_LOCK"; exit 1
fi

log "▶ install-usb 启动 | U 盘: $USB_MOUNT"

# 注入 DNS
printf 'nameserver 223.5.5.5\nnameserver 8.8.8.8\n' > /etc/resolv.conf
restore_dns() {
    printf 'nameserver 127.0.0.1\n' > /etc/resolv.conf
    /etc/init.d/dnsmasq restart 2>/dev/null || true
    log "✅ DNS 已恢复"
}
trap restore_dns EXIT

# apk update
apk update 2>&1 | logger -t "$LOG_TAG"

# I/O 优化
USB_DEV=$(mount | grep " $USB_MOUNT " | awk '{print $1}' | sed 's|/dev/||;s|[0-9]*$||')
SCHED_FILE="/sys/block/${USB_DEV}/queue/scheduler"
if [ -f "$SCHED_FILE" ]; then
    grep -q "mq-deadline" "$SCHED_FILE" 2>/dev/null \
        && echo "mq-deadline" > "$SCHED_FILE" 2>/dev/null \
        || echo "deadline"    > "$SCHED_FILE" 2>/dev/null
fi
echo "10" > /proc/sys/vm/swappiness 2>/dev/null

# Swap
if [ ! -f "$USB_MOUNT/swapfile" ]; then
    log "⏳ 创建 256MB Swap..."
    dd if=/dev/zero of="$USB_MOUNT/swapfile" bs=1M count=256 2>/dev/null \
        && mkswap "$USB_MOUNT/swapfile" >/dev/null 2>&1 \
        && chmod 600 "$USB_MOUNT/swapfile" \
        && swapon "$USB_MOUNT/swapfile" 2>/dev/null \
        && log "✅ Swap 已启用"
else
    swapon "$USB_MOUNT/swapfile" 2>/dev/null || true
fi

INSTALL_FAIL=0
install_pkg() {
    local out ret
    log "📦 安装: $*"
    out=$(apk add --no-cache "$@" 2>&1); ret=$?
    echo "$out" | logger -t "$LOG_TAG"
    [ "$ret" -eq 0 ] && log "✅ 完成: $*" || {
        log "⚠️  失败: $*"
        INSTALL_FAIL=$((INSTALL_FAIL + 1))
    }
}

# Flash 检查
FLASH_FREE=$(df /overlay 2>/dev/null | awk 'NR==2 {print $4}')
SKIP_LUCI=0
[ "${FLASH_FREE:-99999}" -lt 30000 ] && SKIP_LUCI=1

# ── Docker ───────────────────────────────────────────────────
mkdir -p "$USB_MOUNT/docker-data" /etc/docker
cat > /etc/docker/daemon.json << EOF
{
  "data-root": "$USB_MOUNT/docker-data",
  "log-driver": "json-file",
  "log-opts": {"max-size": "5m", "max-file": "2"}
}
EOF

log "📦 强制安装: docker dockerd"
out=$(apk add --no-cache --force-broken-world docker dockerd 2>&1)
ret=$?
echo "$out" | logger -t "$LOG_TAG"
if [ "$ret" -eq 0 ]; then
    log "✅ 完成: docker dockerd"
else
    log "⚠️  失败: docker dockerd"
    INSTALL_FAIL=$((INSTALL_FAIL + 1))
fi

if [ "$SKIP_LUCI" -eq 0 ]; then
    install_pkg luci-app-dockerman luci-i18n-dockerman-zh-cn
else
    log "ℹ️  跳过 luci-app-dockerman（Flash 紧张）"
fi

mountpoint -q /sys/fs/cgroup || mount -t cgroup2 none /sys/fs/cgroup 2>/dev/null || true
/etc/init.d/dockerd enable 2>/dev/null || true
/etc/init.d/dockerd start  2>/dev/null || true
log "✅ Docker 已启动 | 数据: $USB_MOUNT/docker-data"

# ── AdGuardHome ──────────────────────────────────────────────
install_pkg adguardhome
[ "$SKIP_LUCI" -eq 0 ] && install_pkg luci-app-adguardhome \
    || log "ℹ️  跳过 luci-app-adguardhome（Flash 紧张，用 :3000 管理）"

AGH_USB="$USB_MOUNT/AdGuardHome"
mkdir -p "$AGH_USB/config" "$AGH_USB/data"
/etc/init.d/adguardhome stop 2>/dev/null || true

if [ ! -f "$AGH_USB/config/adguardhome.yaml" ] \
&& [ -f "/etc/adguardhome/adguardhome.yaml" ]; then
    cp -a /etc/adguardhome/. "$AGH_USB/config/" 2>/dev/null
    log "✅ AGH 配置迁移至 U 盘（首次）"
else
    log "ℹ️  AGH 配置已存在，复用"
fi

if [ ! "$(ls -A "$AGH_USB/data" 2>/dev/null)" ] \
&& [ -d "/var/lib/adguardhome" ] \
&& [ "$(ls -A /var/lib/adguardhome 2>/dev/null)" ]; then
    cp -a /var/lib/adguardhome/. "$AGH_USB/data/" 2>/dev/null
    log "✅ AGH 数据迁移至 U 盘（首次）"
fi

uci set adguardhome.config.config_file="$AGH_USB/config/adguardhome.yaml" 2>/dev/null || true
uci set adguardhome.config.work_dir="$AGH_USB/data" 2>/dev/null || true
uci commit adguardhome 2>/dev/null || true
/etc/init.d/adguardhome enable 2>/dev/null || true
/etc/init.d/adguardhome start  2>/dev/null || true
log "✅ AdGuardHome 已启动 | $AGH_USB"

# 收尾
rm -f /tmp/luci-indexcache* 2>/dev/null || true
/etc/init.d/rpcd restart 2>/dev/null || true
rm -f "$USB_LOCK"

if [ "$INSTALL_FAIL" -eq 0 ]; then
    touch "$USB_DONE"
    log "🎉 install-usb 完成"
else
    log "⚠️  有 $INSTALL_FAIL 组包失败"
fi
USBEOF
chmod +x /usr/lib/install-extras/install-usb.sh

# ══════════════════════════════════════════════════════════════
# run.sh：基础包 + 条件调用 install-usb.sh
# ══════════════════════════════════════════════════════════════
cat > /usr/lib/install-extras/run.sh << 'RUNEOF'
#!/bin/sh
LOG_TAG="install-extras"
DONE_FLAG="/etc/install-extras-done"
LOCK="/tmp/install-extras.lock"

log() { echo "[$(date '+%H:%M:%S')] $1" | logger -t "$LOG_TAG"; }

[ -f "$DONE_FLAG" ] && exit 0
[ -f "$LOCK" ]      && exit 0
touch "$LOCK"

log "▶ Final v3.4 启动"

for f in /etc/apk/repositories.d/*.list; do
    [ -f "$f" ] && sed -i '/\/video\//d' "$f"
done

# 等待网络
sleep 10
NETWORK_OK=0
for i in $(seq 1 36); do
    ping -c1 -W2 223.5.5.5 >/dev/null 2>&1 && { NETWORK_OK=1; break; }
    ping -c1 -W2 8.8.8.8  >/dev/null 2>&1 && { NETWORK_OK=1; break; }
    sleep 5
done
if [ "$NETWORK_OK" -eq 0 ]; then
    log "❌ 网络不可达，等待下次重试"
    rm -f "$LOCK"; exit 1
fi

printf 'nameserver 223.5.5.5\nnameserver 8.8.8.8\n' > /etc/resolv.conf
log "✅ 网络就绪，DNS 已注入"

restore_dns() {
    printf 'nameserver 127.0.0.1\n' > /etc/resolv.conf
    /etc/init.d/dnsmasq restart 2>/dev/null || true
    log "✅ DNS 已恢复"
}
trap restore_dns EXIT

# apk update（先做，顺带给 automount 留时间）
apk_update_ok() {
    local out n
    out=$(apk update 2>&1); echo "$out" | logger -t "$LOG_TAG"
    n=$(echo "$out" | grep -oE 'OK: [0-9]+ distinct' | grep -oE '[0-9]+')
    [ "${n:-0}" -gt 800 ] && return 0
    log "⚠️  HTTPS 失败（${n:-0}包），切换 HTTP..."
    for f in /etc/apk/repositories.d/*.list; do
        [ -f "$f" ] && sed -i 's|https://|http://|g' "$f"
    done
    out=$(apk update 2>&1); echo "$out" | logger -t "$LOG_TAG"
    n=$(echo "$out" | grep -oE 'OK: [0-9]+ distinct' | grep -oE '[0-9]+')
    [ "${n:-0}" -gt 800 ] && return 0
    return 1
}

if ! apk_update_ok; then
    log "❌ apk update 失败，等待下次重试"
    rm -f "$LOCK"; exit 1
fi
log "✅ apk update 成功"

# U 盘检测（apk update 后检测，automount 已有足够时间）
USB_MOUNT=""
log "⏳ 检测 U 盘..."
for i in $(seq 1 5); do
    for mp in /mnt/sda1 /mnt/sda /mnt/usb /mnt/mmcblk0p1; do
        mountpoint -q "$mp" 2>/dev/null && { USB_MOUNT="$mp"; break 2; }
    done
    sleep 2
done
[ -n "$USB_MOUNT" ] && log "✅ U 盘: $USB_MOUNT" || log "⚠️  无 U 盘"

INSTALL_FAIL=0
install_pkg() {
    local out ret
    log "📦 安装: $*"
    out=$(apk add --no-cache "$@" 2>&1); ret=$?
    echo "$out" | logger -t "$LOG_TAG"
    [ "$ret" -eq 0 ] && log "✅ 完成: $*" || {
        log "⚠️  失败: $*"
        INSTALL_FAIL=$((INSTALL_FAIL + 1))
    }
}

FLASH_FREE=$(df /overlay 2>/dev/null | awk 'NR==2 {print $4}')
FLASH_USED=$(df /overlay 2>/dev/null | awk 'NR==2 {print $5}' | tr -d '%')
log "📊 Flash: 已用 ${FLASH_USED:-?}% | 剩余 ${FLASH_FREE:-?}KB"
SKIP_LUCI=0
[ "${FLASH_FREE:-99999}" -lt 30000 ] && SKIP_LUCI=1

# 基础工具
install_pkg smartmontools bash htop fdisk lsblk

# Samba
install_pkg samba4-server
[ "$SKIP_LUCI" -eq 0 ] && install_pkg luci-app-samba4 luci-i18n-samba4-zh-cn \
    || log "ℹ️  跳过 luci-app-samba4（Flash 紧张）"

/etc/init.d/dbus restart 2>/dev/null || true
sleep 2
/etc/init.d/avahi-daemon restart 2>/dev/null || true
log "✅ dbus / avahi 已重载"

# 重型应用
if [ -n "$USB_MOUNT" ] && [ ! -f /etc/install-usb-done ]; then
    /usr/lib/install-extras/install-usb.sh "$USB_MOUNT"
elif [ -z "$USB_MOUNT" ]; then
    log "ℹ️  无 U 盘，跳过重型应用"
    log "ℹ️  插入 U 盘后会自动触发安装，或执行: rm -f /etc/install-usb-done && /usr/lib/install-extras/install-usb.sh"
fi

rm -f /tmp/luci-indexcache* 2>/dev/null || true
/etc/init.d/rpcd restart 2>/dev/null || true
rm -f /etc/hotplug.d/iface/99-install-extras
rm -f "$LOCK"

if [ "$INSTALL_FAIL" -eq 0 ]; then
    FLASH_FINAL=$(df /overlay 2>/dev/null | awk 'NR==2 {print $5}' | tr -d '%')
    touch "$DONE_FLAG"
    log "🎉 Final v3.4 完成！Flash: ${FLASH_FINAL:-?}%"
else
    log "⚠️  有 $INSTALL_FAIL 组包失败，注册重试"
    mkdir -p /etc/hotplug.d/iface
    cat > /etc/hotplug.d/iface/99-install-extras << 'RETRYEOF'
#!/bin/sh
[ "$ACTION" = "ifup" ]          || exit 0
[ -f /etc/install-extras-done ] && exit 0
case "$INTERFACE" in wan|wan6|pppoe-wan) ;; *) exit 0 ;; esac
/usr/lib/install-extras/run.sh &
RETRYEOF
    chmod +x /etc/hotplug.d/iface/99-install-extras
fi
RUNEOF
chmod +x /usr/lib/install-extras/run.sh

# ── WAN UP 触发器 ────────────────────────────────────────────────────
mkdir -p /etc/hotplug.d/iface
cat > /etc/hotplug.d/iface/99-install-extras << 'HOTEOF'
#!/bin/sh
[ "$ACTION" = "ifup" ]          || exit 0
[ -f /etc/install-extras-done ] && exit 0
case "$INTERFACE" in wan|wan6|pppoe-wan) ;; *) exit 0 ;; esac
/usr/lib/install-extras/run.sh &
HOTEOF
chmod +x /etc/hotplug.d/iface/99-install-extras

# ── USB 热插触发器 ───────────────────────────────────────────────────
mkdir -p /etc/hotplug.d/block
cat > /etc/hotplug.d/block/99-usb-extras << 'BLOCKEOF'
#!/bin/sh
[ "$ACTION" = "add" ] || exit 0

# 等待 automount 完成
sleep 5

USB_MOUNT=""
for mp in /mnt/sda1 /mnt/sda /mnt/usb /mnt/mmcblk0p1; do
    mountpoint -q "$mp" 2>/dev/null && { USB_MOUNT="$mp"; break; }
done
[ -n "$USB_MOUNT" ] || exit 0

# 重型应用已安装则跳过
[ -f /etc/install-usb-done ] && exit 0

# WAN 必须在线
ping -c1 -W3 223.5.5.5 >/dev/null 2>&1 || \
ping -c1 -W3 8.8.8.8  >/dev/null 2>&1 || exit 0

logger -t "usb-extras" "✅ U 盘热插检测: $USB_MOUNT，触发重型应用安装"
/usr/lib/install-extras/install-usb.sh "$USB_MOUNT" &
BLOCKEOF
chmod +x /etc/hotplug.d/block/99-usb-extras

logger -t "install-extras" "✅ Final v3.4 已安装"
exit 0
