# files/etc/uci-defaults/96-install-extras
#!/bin/sh
# 96-install-extras (128M 极简融合版)
# 设计目标：
# - 不动系统（不自动扩容）
# - 自动识别 overlay / U盘
# - 所有重数据外移
# - 超低占用 / 高稳定

DONE_FLAG="/etc/install-extras-done"
[ -f "$DONE_FLAG" ] && exit 0

mkdir -p /usr/lib/install-extras

cat > /usr/lib/install-extras/run.sh << 'RUNEOF'
#!/bin/sh
LOG_TAG="install-extras-mini"
DONE_FLAG="/etc/install-extras-done"
LOCK="/tmp/install-extras.lock"

log() { echo "[$(date '+%H:%M:%S')] $1" | logger -t "$LOG_TAG"; }

[ -f "$DONE_FLAG" ] && exit 0
[ -f "$LOCK" ] && exit 0
touch "$LOCK"

log "▶ 极简融合版启动"

# ── 网络检测 ──
sleep 5
for i in $(seq 1 15); do
    ping -c1 -W1 8.8.8.8 >/dev/null 2>&1 && break
    sleep 2
done

# ── DNS 修复 ──
echo "nameserver 8.8.8.8" > /etc/resolv.conf

# ── 判断 overlay ──
if mount | grep -q "on /overlay type overlay"; then
    OVERLAY=1
    DATA_ROOT="/overlay"
    log "✅ overlay 模式"
else
    OVERLAY=0
    log "ℹ️ 非 overlay 模式"
fi

# ── 检测 U盘（仅 fallback） ──
USB_MOUNT=""
for mp in /mnt/sda1 /mnt/sda /mnt/usb; do
    mountpoint -q "$mp" 2>/dev/null && USB_MOUNT="$mp" && break
done

if [ "$OVERLAY" -eq 0 ] && [ -n "$USB_MOUNT" ]; then
    DATA_ROOT="$USB_MOUNT"
    log "📌 使用U盘: $USB_MOUNT"
fi

if [ -z "$DATA_ROOT" ]; then
    log "⚠️ 无 overlay 且无U盘，跳过重型应用"
fi

# ── 更新源 ──
apk update || {
    log "❌ apk update失败"
    rm -f "$LOCK"
    exit 1
}

# ── 基础工具 ──
apk add --no-cache bash htop >/dev/null 2>&1 && log "✅ 基础工具完成"

# ── 重型应用（仅有存储时） ──
if [ -n "$DATA_ROOT" ]; then

    # Docker
    mkdir -p "$DATA_ROOT/docker-data" /etc/docker
    echo '{"data-root":"'$DATA_ROOT'/docker-data"}' > /etc/docker/daemon.json

    apk add --no-cache docker dockerd >/dev/null 2>&1 && log "✅ Docker安装完成"

    mountpoint -q /sys/fs/cgroup || mount -t cgroup2 none /sys/fs/cgroup 2>/dev/null
    /etc/init.d/dockerd enable 2>/dev/null
    /etc/init.d/dockerd start 2>/dev/null

    # AdGuardHome
    apk add --no-cache adguardhome >/dev/null 2>&1 && log "✅ AGH安装完成"

    AGH_DIR="$DATA_ROOT/adguardhome"
    mkdir -p "$AGH_DIR/data" "$AGH_DIR/config"

    uci set adguardhome.config.work_dir="$AGH_DIR/data" 2>/dev/null
    uci set adguardhome.config.config_file="$AGH_DIR/config/adguardhome.yaml" 2>/dev/null
    uci commit adguardhome 2>/dev/null

    /etc/init.d/adguardhome enable 2>/dev/null
    /etc/init.d/adguardhome start 2>/dev/null

    log "✅ AGH已运行（数据在 $DATA_ROOT）"
fi

rm -f "$LOCK"
touch "$DONE_FLAG"
log "🎉 完成（极简版）"
RUNEOF

chmod +x /usr/lib/install-extras/run.sh

# ── hotplug ──
mkdir -p /etc/hotplug.d/iface
cat > /etc/hotplug.d/iface/99-install-extras << 'EOF'
#!/bin/sh
[ "$ACTION" = "ifup" ] || exit 0
/usr/lib/install-extras/run.sh &
EOF

chmod +x /etc/hotplug.d/iface/99-install-extras

logger -t "install-extras-mini" "✅ 极简脚本已安装"
exit 0
