# files/etc/uci-defaults/96-install-extras
#!/bin/sh
# 96-install-extras
# hotplug iface 触发：WAN UP 后自动安装软件包

DONE_FLAG="/etc/install-extras-done"
[ -f "$DONE_FLAG" ] && exit 0

mkdir -p /usr/lib/install-extras

cat > /usr/lib/install-extras/run.sh << 'RUNEOF'
#!/bin/sh
LOG_TAG="install-extras"
DONE_FLAG="/etc/install-extras-done"
LOCK="/tmp/install-extras.lock"

log() { echo "[$(date '+%H:%M:%S')] $1" | logger -t "$LOG_TAG"; }

[ -f "$DONE_FLAG" ] && exit 0
[ -f "$LOCK" ] && exit 0
touch "$LOCK"

log "▶ install-extras 开始"

# ── 清理 video feed ──────────────────────────────────────────────────
for f in /etc/apk/repositories.d/*.list; do
    [ -f "$f" ] && sed -i '/\/video\//d' "$f"
done

# ── 等待网络 ─────────────────────────────────────────────────────────
sleep 10
NETWORK_OK=0
for i in $(seq 1 60); do
    ping -c 1 -W 2 223.5.5.5 >/dev/null 2>&1 && { NETWORK_OK=1; break; }
    ping -c 1 -W 2 8.8.8.8  >/dev/null 2>&1 && { NETWORK_OK=1; break; }
    sleep 3
done
if [ "$NETWORK_OK" -eq 0 ]; then
    log "❌ 网络不可达，等下次 WAN UP 重试"
    rm -f "$LOCK"; exit 1
fi
log "✅ 网络就绪"

# ── 强制写入公共 DNS，绕过 smartdns 时序问题 ─────────────────────────
printf 'nameserver 223.5.5.5\nnameserver 8.8.8.8\n' > /etc/resolv.conf
log "✅ DNS 注入：223.5.5.5 / 8.8.8.8"

# ── apk update ───────────────────────────────────────────────────────
apk_update_ok() {
    local out n
    out=$(apk update 2>&1)
    echo "$out" | logger -t "$LOG_TAG"
    n=$(echo "$out" | grep -oE 'OK: [0-9]+ distinct' | grep -oE '[0-9]+')
    [ "${n:-0}" -gt 1000 ] && return 0

    log "⚠️  HTTPS 包数不足（${n:-0}），切换 HTTP..."
    for f in /etc/apk/repositories.d/*.list; do
        [ -f "$f" ] && sed -i \
            's|https://downloads.immortalwrt.org|http://downloads.immortalwrt.org|g' "$f"
    done
    out=$(apk update 2>&1)
    echo "$out" | logger -t "$LOG_TAG"
    n=$(echo "$out" | grep -oE 'OK: [0-9]+ distinct' | grep -oE '[0-9]+')
    [ "${n:-0}" -gt 1000 ] && return 0
    return 1
}

if ! apk_update_ok; then
    log "❌ apk update 失败，等下次重试"
    rm -f "$LOCK"; exit 1
fi
log "✅ apk update 成功"

# ── 安装函数 ─────────────────────────────────────────────────────────
INSTALL_FAIL=0
install_pkg() {
    log "📦 安装: $*"
    if apk add --no-cache "$@" 2>&1 | logger -t "$LOG_TAG"; then
        log "✅ 完成: $*"
    else
        log "⚠️  失败: $*"
        INSTALL_FAIL=$((INSTALL_FAIL + 1))
    fi
}

# docker/dockerman 已编译内置，不在此安装
# lucky/fileassistant/crontabs 官方源不存在，已移除
install_pkg adguardhome luci-app-adguardhome
install_pkg samba4-server luci-app-samba4 luci-i18n-samba4-zh-cn
install_pkg smartmontools bash htop fdisk lsblk

# ── 配置 Docker 数据目录 ─────────────────────────────────────────────
if command -v dockerd >/dev/null 2>&1; then
    mkdir -p /etc/docker
    USB_MOUNT=""
    for mp in /mnt/sda1 /mnt/sda /mnt/usb /mnt/mmcblk0p1; do
        mountpoint -q "$mp" 2>/dev/null && { USB_MOUNT="$mp"; break; }
    done
    if [ -n "$USB_MOUNT" ]; then
        mkdir -p "$USB_MOUNT/docker-data"
        printf '{\n  "data-root": "%s/docker-data",\n  "log-driver": "json-file",\n  "log-opts": {"max-size": "10m", "max-file": "3"}\n}\n' \
            "$USB_MOUNT" > /etc/docker/daemon.json
        log "✅ Docker 数据目录: $USB_MOUNT/docker-data"
    else
        printf '{\n  "log-driver": "json-file",\n  "log-opts": {"max-size": "5m", "max-file": "2"}\n}\n' \
            > /etc/docker/daemon.json
        log "⚠️  未检测到 U 盘，Docker 数据存储在内置 Flash"
    fi
    mountpoint -q /sys/fs/cgroup || mount -t cgroup2 none /sys/fs/cgroup 2>/dev/null || true
    /etc/init.d/dockerd enable 2>/dev/null || true
    /etc/init.d/dockerd start  2>/dev/null || true
fi

# ── 刷新 LuCI ────────────────────────────────────────────────────────
rm -f /tmp/luci-indexcache* 2>/dev/null
/etc/init.d/rpcd restart 2>/dev/null || true

rm -f /etc/hotplug.d/iface/99-install-extras
rm -f "$LOCK"

# ── 写 DONE_FLAG ─────────────────────────────────────────────────────
if [ "$INSTALL_FAIL" -eq 0 ]; then
    touch "$DONE_FLAG"
    log "🎉 install-extras 全部完成"
else
    log "⚠️  有 $INSTALL_FAIL 组包失败，下次 WAN UP 重试"
    mkdir -p /etc/hotplug.d/iface
    cat > /etc/hotplug.d/iface/99-install-extras << 'HOTPLUGEOF'
#!/bin/sh
[ "$ACTION" = "ifup" ] || exit 0
[ -f /etc/install-extras-done ] && exit 0
case "$INTERFACE" in
    wan|wan6|pppoe-wan) ;;
    *) exit 0 ;;
esac
/usr/lib/install-extras/run.sh &
HOTPLUGEOF
    chmod +x /etc/hotplug.d/iface/99-install-extras
fi
RUNEOF
chmod +x /usr/lib/install-extras/run.sh

# ── 安装 hotplug 触发器 ──────────────────────────────────────────────
mkdir -p /etc/hotplug.d/iface
cat > /etc/hotplug.d/iface/99-install-extras << 'HOTPLUGEOF'
#!/bin/sh
[ "$ACTION" = "ifup" ] || exit 0
[ -f /etc/install-extras-done ] && exit 0
case "$INTERFACE" in
    wan|wan6|pppoe-wan) ;;
    *) exit 0 ;;
esac
/usr/lib/install-extras/run.sh &
HOTPLUGEOF
chmod +x /etc/hotplug.d/iface/99-install-extras

logger -t "96-install-extras" "✅ hotplug 触发器已安装"
exit 0
