# files/etc/uci-defaults/96-install-extras
#!/bin/sh
# 96-install-extras
# 首次启动后台安装大体积包，不阻塞启动流程
LOG_TAG="install-extras"
DONE_FLAG="/etc/install-extras-done"

[ -f "$DONE_FLAG" ] && exit 0

nohup sh -c '
LOG_TAG="install-extras"
DONE_FLAG="/etc/install-extras-done"
[ -f "$DONE_FLAG" ] && exit 0

log() { echo "[$(date +%H:%M:%S)] $1" | logger -t "$LOG_TAG"; }

log "⏳ 等待系统启动完成..."
sleep 30

# 修复 video feed
REPO_DIR="/etc/apk/repositories.d"
for f in "$REPO_DIR"/*.list; do
    [ -f "$f" ] && sed -i "/\/video\//d" "$f"
done

# 等待 WAN 网络就绪（最多等 3 分钟）
log "⏳ 等待 WAN 网络..."
for i in $(seq 1 90); do
    ping -c 1 -W 2 223.5.5.5 >/dev/null 2>&1 && break
    sleep 2
done

# 更新索引
apk update 2>&1 | logger -t "$LOG_TAG"
true

# 安装各包（逐个安装，失败不影响其他包）
for pkg in "dockerd docker docker-compose cgroupfs-mount" "samba4-server" "adguardhome" "lucky luci-app-lucky luci-i18n-lucky-zh-cn" "luci-app-fileassistant"; do
    log "📦 安装: $pkg"
    apk add --no-cache $pkg 2>&1 | logger -t "$LOG_TAG" \
        && log "✅ $pkg 安装完成" \
        || log "⚠️  $pkg 安装失败"
done

# 刷新 LuCI
rm -f /tmp/luci-indexcache* 2>/dev/null
/etc/init.d/rpcd restart 2>/dev/null || true

touch "$DONE_FLAG"
log "🎉 install-extras 完成"
' >> /tmp/install-extras.log 2>&1 &

exit 0
