# files/etc/uci-defaults/96-install-extras
#!/bin/sh
# 96-install-extras
# WAN UP 后自动安装软件包，数据外移至 U 盘

DONE_FLAG="/etc/install-extras-done"
[ -f "$DONE_FLAG" ] && exit 0

mkdir -p /usr/lib/install-extras

cat > /usr/lib/install-extras/run.sh << 'RUNEOF'
#!/bin/sh
LOG_TAG="install-extras"
DONE_FLAG="/etc/install-extras-done"
LOCK="/tmp/install-extras.lock"

log() { echo "[$(date '+%H:%M:%S')] $1" | logger -t "$LOG_TAG"; }

[ -f "$DONE_FLAG" ] && exit 0
[ -f "$LOCK" ] && exit 0
touch "$LOCK"

log "▶ install-extras 开始"

# ── 等待网络 ─────────────────────────────────────────────────────────
sleep 10
NETWORK_OK=0
for i in $(seq 1 60); do
    ping -c 1 -W 2 223.5.5.5 >/dev/null 2>&1 && { NETWORK_OK=1; break; }
    ping -c 1 -W 2 8.8.8.8  >/dev/null 2>&1 && { NETWORK_OK=1; break; }
    sleep 3
done
if [ "$NETWORK_OK" -eq 0 ]; then
    log "❌ 网络不可达，等下次 WAN UP 重试"
    rm -f "$LOCK"; exit 1
fi
log "✅ 网络就绪"

# ── DNS 修复 ─────────────────────────────────────────────────────────
printf 'nameserver 223.5.5.5\nnameserver 8.8.8.8\n' > /etc/resolv.conf
log "✅ DNS 注入"

# ── 检测 U 盘（等待最多 30 秒）───────────────────────────────────────
USB_MOUNT=""
log "⏳ 检测 U 盘..."
for i in $(seq 1 15); do
    for mp in /mnt/sda1 /mnt/sda /mnt/usb /mnt/mmcblk0p1; do
        mountpoint -q "$mp" 2>/dev/null && { USB_MOUNT="$mp"; break 2; }
    done
    sleep 2
done
if [ -n "$USB_MOUNT" ]; then
    log "✅ U 盘已挂载: $USB_MOUNT"
else
    log "⚠️  未检测到 U 盘，重型应用将跳过"
fi

# ── apk update（带 HTTPS→HTTP 降级重试）──────────────────────────────
apk_update_ok() {
    local out n
    out=$(apk update 2>&1)
    echo "$out" | logger -t "$LOG_TAG"
    n=$(echo "$out" | grep -oE 'OK: [0-9]+ distinct' | grep -oE '[0-9]+')
    [ "${n:-0}" -gt 1000 ] && return 0

    log "⚠️  HTTPS 包数不足（${n:-0}），切换 HTTP..."
    for f in /etc/apk/repositories.d/*.list; do
        [ -f "$f" ] && sed -i \
            's|https://downloads.immortalwrt.org|http://downloads.immortalwrt.org|g' "$f"
    done
    out=$(apk update 2>&1)
    echo "$out" | logger -t "$LOG_TAG"
    n=$(echo "$out" | grep -oE 'OK: [0-9]+ distinct' | grep -oE '[0-9]+')
    [ "${n:-0}" -gt 1000 ] && return 0
    return 1
}

if ! apk_update_ok; then
    log "❌ apk update 失败，等下次重试"
    rm -f "$LOCK"; exit 1
fi
log "✅ apk update 成功"

# ── 安装函数（独立计数）──────────────────────────────────────────────
INSTALL_FAIL=0
install_pkg() {
    local out ret
    log "📦 安装: $*"
    out=$(apk add --no-cache "$@" 2>&1)
    ret=$?
    echo "$out" | logger -t "$LOG_TAG"
    [ "$ret" -eq 0 ] && log "✅ 完成: $*" || {
        log "⚠️  失败: $*"
        INSTALL_FAIL=$((INSTALL_FAIL + 1))
    }
}

# ── 轻量工具：无论有无 U 盘直接装 ───────────────────────────────────
install_pkg bash htop smartmontools fdisk lsblk

# ── Samba：装 Flash，装完修复 dbus/avahi 时序问题 ────────────────────
install_pkg samba4-server luci-app-samba4 luci-i18n-samba4-zh-cn
/etc/init.d/dbus restart 2>/dev/null || true
sleep 2
/etc/init.d/avahi-daemon restart 2>/dev/null || true
log "✅ dbus 已重载，avahi-daemon 已启动"

# ── 重型应用：强依赖 U 盘，无盘跳过 ─────────────────────────────────
if [ -n "$USB_MOUNT" ]; then
    log "⚙️  U 盘就绪，开始安装重型应用..."

    # Swap（第一件事，缓解安装时内存压力）
    if [ -f "$USB_MOUNT/swapfile" ]; then
        swapon "$USB_MOUNT/swapfile" 2>/dev/null && \
            log "✅ Swap 已开启" || \
            log "⚠️  Swap 开启失败"
    else
        log "ℹ️  如需 Swap 可执行:"
        log "ℹ️  dd if=/dev/zero of=$USB_MOUNT/swapfile bs=1M count=512 && mkswap $USB_MOUNT/swapfile"
    fi

    # Docker（先装，体积最大优先确认 Flash 空间）
    mkdir -p "$USB_MOUNT/docker-data" /etc/docker
    printf '{\n  "data-root": "%s/docker-data",\n  "log-driver": "json-file",\n  "log-opts": {"max-size": "10m", "max-file": "3"}\n}\n' \
        "$USB_MOUNT" > /etc/docker/daemon.json
    log "✅ daemon.json 已写入"

    install_pkg docker dockerd
    install_pkg luci-app-dockerman luci-i18n-dockerman-zh-cn

    mountpoint -q /sys/fs/cgroup || \
        mount -t cgroup2 none /sys/fs/cgroup 2>/dev/null || true
    /etc/init.d/dockerd enable 2>/dev/null || true
    /etc/init.d/dockerd start  2>/dev/null || true
    log "✅ Docker 已启动，数据 → $USB_MOUNT/docker-data"

    # AdGuardHome（Docker 装完后装，Flash 空间已确认够用）
    install_pkg adguardhome luci-app-adguardhome

    AGH_USB="$USB_MOUNT/AdGuardHome"
    mkdir -p "$AGH_USB/data" "$AGH_USB/config"
    /etc/init.d/adguardhome stop 2>/dev/null || true

    # 优先保留 U 盘已有配置，重刷不覆盖
    if [ ! -f "$AGH_USB/config/adguardhome.yaml" ] && \
       [ -f "/etc/adguardhome/adguardhome.yaml" ]; then
        cp -a /etc/adguardhome/. "$AGH_USB/config/"
        log "✅ 新环境，初始配置已写入 U 盘"
    else
        log "ℹ️  U 盘已有旧配置，直接复用"
    fi

    uci set adguardhome.config.config_file="$AGH_USB/config/adguardhome.yaml"
    uci set adguardhome.config.work_dir="$AGH_USB/data"
    uci commit adguardhome
    /etc/init.d/adguardhome enable 2>/dev/null || true
    /etc/init.d/adguardhome start  2>/dev/null || true
    log "✅ AdGuardHome 已启动，数据 → $AGH_USB"

else
    log "ℹ️  无 U 盘，跳过 Docker 和 AdGuardHome"
    log "ℹ️  插入 U 盘后执行: rm /etc/install-extras-done && reboot"
fi

# ── 刷新 LuCI ────────────────────────────────────────────────────────
rm -f /tmp/luci-indexcache* 2>/dev/null
/etc/init.d/rpcd restart 2>/dev/null || true

rm -f "$LOCK"

# ── 写 DONE_FLAG（仅全部成功时）──────────────────────────────────────
if [ "$INSTALL_FAIL" -eq 0 ]; then
    touch "$DONE_FLAG"
    log "🎉 install-extras 全部完成"
else
    log "⚠️  有 $INSTALL_FAIL 组包失败，下次 WAN UP 自动重试"
    # 重新安装 hotplug 触发器等待下次重试
    mkdir -p /etc/hotplug.d/iface
    cat > /etc/hotplug.d/iface/99-install-extras << 'HOTPLUGEOF'
#!/bin/sh
[ "$ACTION" = "ifup" ] || exit 0
[ -f /etc/install-extras-done ] && exit 0
case "$INTERFACE" in
    wan|wan6|pppoe-wan) ;;
    *) exit 0 ;;
esac
/usr/lib/install-extras/run.sh &
HOTPLUGEOF
    chmod +x /etc/hotplug.d/iface/99-install-extras
fi
RUNEOF
chmod +x /usr/lib/install-extras/run.sh

# ── hotplug 触发器（仅 WAN 接口）────────────────────────────────────
mkdir -p /etc/hotplug.d/iface
cat > /etc/hotplug.d/iface/99-install-extras << 'HOTPLUGEOF'
#!/bin/sh
[ "$ACTION" = "ifup" ] || exit 0
[ -f /etc/install-extras-done ] && exit 0
case "$INTERFACE" in
    wan|wan6|pppoe-wan) ;;
    *) exit 0 ;;
esac
/usr/lib/install-extras/run.sh &
HOTPLUGEOF
chmod +x /etc/hotplug.d/iface/99-install-extras

logger -t "install-extras" "✅ hotplug 触发器已安装"
exit 0
