# files/etc/uci-defaults/97-install-docker
#!/bin/sh
# 97-install-docker
# dockerd 由 96-install-extras 后台安装
# 本脚本只负责在系统启动后确保 dockerd 服务自启
# 因为 96 是后台执行的，这里不需要等待

# 把 dockerd 启动也放到后台，等 96 安装完再启动
nohup sh -c '
LOG_TAG="start-docker"
log() { echo "[$(date +%H:%M:%S)] $1" | logger -t "$LOG_TAG"; }

# 等待 dockerd 二进制出现（最多等 5 分钟）
for i in $(seq 1 150); do
    command -v dockerd >/dev/null 2>&1 && break
    sleep 2
done

if ! command -v dockerd >/dev/null 2>&1; then
    log "❌ dockerd 未安装，跳过启动"
    exit 0
fi

# 挂载 cgroup2
mount | grep -q cgroup2 || mount -t cgroup2 none /sys/fs/cgroup 2>/dev/null || true

# 启用并启动
service dockerd enable 2>/dev/null
service dockerd start  2>/dev/null
sleep 3
service dockerd status 2>/dev/null | grep -q running \
    && log "✅ dockerd 启动成功" \
    || log "⚠️  dockerd 未运行，请手动: service dockerd start"
' >> /tmp/install-extras.log 2>&1 &

exit 0
