# files/etc/uci-defaults/97-install-docker
#!/bin/sh
# 99-custom-ip - 设置 LAN IP 为 192.168.2.1
uci set network.lan.ipaddr='192.168.2.1'
uci set network.lan.netmask='255.255.255.0'
uci commit network
exit 0
