# files/etc/uci-defaults/98-fix-apk-feeds
#!/bin/sh
# 98-fix-apk-feeds
# 修复 immortalwrt snapshot 中 video feed 不存在导致的 apk update 报错
# 执行时机：首次启动，uci-defaults 自动运行后自删

REPO_LIST_DIR="/etc/apk/repositories.d"

# 等待网络就绪（最多等 15 秒）
for i in $(seq 1 15); do
    ping -c 1 -W 2 8.8.8.8 >/dev/null 2>&1 && break
    sleep 1
done

# 删除不存在的 video feed
for f in "$REPO_LIST_DIR"/*.list; do
    [ -f "$f" ] || continue
    sed -i '/\/video\//d' "$f"
done

echo "✅ 98-fix-apk-feeds: video feed 已清理"

# 更新软件源索引
apk update 2>&1 | logger -t apk-fix
echo "✅ 98-fix-apk-feeds: apk update 完成"

exit 0
