# files/etc/uci-defaults/98-fix-apk-feeds
#!/bin/sh
# 98-fix-apk-feeds
# 清理 apk video feed（会导致 apk update 报错），立即执行不等待网络
REPO_LIST_DIR="/etc/apk/repositories.d"
for f in "$REPO_LIST_DIR"/*.list; do
    [ -f "$f" ] || continue
    sed -i '/\/video\//d' "$f"
done
logger -t "98-fix-apk-feeds" "✅ video feed 已清理"
exit 0
